# pmm_client.ini
[program:pmm-client]
command = pmm-admin add mysql --host db
stdout_logfile = /var/log/pmm-client.log
stderr_logfile = /var/log/pmm-client.log
autostart = true
startretries = 1000000
autorestart = unexpected