# pmm_client.ini
[program:pmm-client]
command = pmm-admin add mysql --host %(ENV_MYSQL_HOST)s
stdout_logfile = /var/log/pmm-client.log
stderr_logfile = /var/log/pmm-client.log
autostart = true
startretries = 1000000
autorestart = unexpected
